# confd/confd.service
[Unit]
Description=Confd
After=haproxy.service

[Service]
ExecStart=/usr/local/bin/confd
Restart=always

[Install]
WantedBy=basic.target